# retro-docker/task-4/gitlab-runner-hub.timer
[Unit]
Description=Gitlab Runner Hub Timer


[Timer]
OnCalendar= *-*-* *:*:*
Unit=gitlab-runner-hub.service


[Install]
WantedBy=default.target