# retro-docker/task-2/docker-viper.timer
[Unit]
Description=Docker Viper Timer


[Timer]
OnCalendar= *-*-* *:*:*
Unit=docker-viper.service


[Install]
WantedBy=default.target